# config/crontab
# if successful add a cron twice a day randomly
* * * * * /app/config/update_certs.sh >> /app/logs/renewal.log 2>&1
